# libldm/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/libldm
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/ldmtool
drwxr-xr-x	root/root	usr/include/
-rw-r--r--	root/root	usr/include/ldm.h
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libldm-1.0.a
-rwxr-xr-x	root/root	usr/lib/libldm-1.0.la
lrwxrwxrwx	root/root	usr/lib/libldm-1.0.so -> libldm-1.0.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libldm-1.0.so.0 -> libldm-1.0.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libldm-1.0.so.0.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/ldm-1.0.pc
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/ldmtool.1.gz
